# roles/kibana/files/kibana.ini
[program:kibana]
user=user
command=/opt/kibana/bin/kibana
autostart=true
autorestart=true
stdout_logfile=/opt/log/kibana.out.log
stderr_logfile=/opt/log/kibana.err.log