# roles/filebeat/files/filebeat.ini
[program:filebeat]
user=root
command=/opt/filebeat/filebeat -e -c filebeat.yml
directory=/opt/filebeat/
autostart=true
autorestart=true
stdout_logfile=/opt/log/filebeat.out.log
stderr_logfile=/opt/log/filebeat.err.log