# roles/logstash/files/logstash.ini
[program:logstash]
user=user
command=/opt/logstash/bin/logstash -f /opt/logstash/config/logstash.conf
autostart=true
autorestart=true
stdout_logfile=/opt/log/logstash.out.log
stderr_logfile=/opt/log/logstash.err.log